# hapro/rootfs/etc/s6-overlay/s6-rc.d/rathole/run
#!/usr/bin/with-contenv bashio
source /usr/src/config/vars.sh
source /usr/src/request.sh

doSupervisorRequest "core/info"
haIp=$(echo $response | jq -r '.data.ip_address')
sed -i "s/{{.PORT}}/80/g" /usr/bin/client.toml
sed -i "s/{{.HOST}}/$haIp/g" /usr/bin/client.toml

bashio::log.info "RATHOLE: Starting rathole"
cat /usr/bin/client.toml
exec /usr/bin/rathole -c /usr/bin/client.toml