# ems/rootfs/etc/services.d/rathole/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# s6-overlay docs: https://github.com/just-containers/s6-overlay

# Add your code here

# Declare variables
declare token

token=$(bashio::config 'token')

curl "https://admin.ems.bitfox.nl/homes/$token/config.toml" -o /usr/bin/client.toml

bashio::log.info "client.toml:"
cat /usr/bin/client.toml
bashio::log.info "token: $token"

exec /usr/bin/rathole -c /usr/bin/client.toml
