# hapro/rootfs/etc/s6-overlay/s6-rc.d/webapi/finish
#!/usr/bin/with-contenv bashio
declare APP_EXIT_CODE=${1}
declare container="webapi"

if [[ "${APP_EXIT_CODE}" -ne 0 ]] && [[ "${APP_EXIT_CODE}" -ne 256 ]] && [[ "${APP_EXIT_CODE}" -ne 137 ]] ; then
  bashio::log.warning "Halt $container with exit code ${APP_EXIT_CODE}"
  echo "${APP_EXIT_CODE}" > /run/s6-linux-init-container-results/exitcode
  exec /run/s6/basedir/bin/halt
elif [[ "${APP_EXIT_CODE}" -eq 137 ]] ; then
  bashio::log.warning "Container $container was killed by OOM killer, restarting"
else
  bashio::log.info "Automatically restarting $container with exit code ${APP_EXIT_CODE}"
fi