# hapro/rootfs/etc/s6-overlay/s6-rc.d/webapi/run
#!/usr/bin/with-contenv bashio
bashio::log.info "WEBAPI: Starting Web API"
exec bun run /usr/src/webapi.ts
