# ems/rootfs/etc/services.d/add_addon/finish
#!/usr/bin/env bashio
# ==============================================================================
# Take down the S6 supervision tree when example fails
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# bashio::log.info "Halt!"
# exec /run/s6/basedir/bin/halt
