# ems/rootfs/etc/services.d/add_addon/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Running CP JOB"

ls /homeassistant

mkdir -p /homeassistant/custom_components

cp -R /usr/src/hass-auth-header/custom_components/* /homeassistant/custom_components/

bashio::log.info "Done CP JOB"

sleep 60000000