# ems/rootfs/etc/services.d/add_addon/run
#!/command/with-contenv bashio
# shellcheck shell=bash
declare token

token=$(bashio::config 'token')
# exec > >(tee -a /homeassistant/runlog.log) 2>&1

bashio::log.info "Running CP JOB"

mkdir -p /homeassistant/custom_components

# Copy only if the directory is not already present
if [ ! -d "/homeassistant/custom_components/auth_header" ]; then
    cp -R /usr/src/hass-auth-header/custom_components/* /homeassistant/custom_components/
fi

# Copy ems.yaml only if it doesn't already exist
if [ ! -f "/homeassistant/ems.yaml" ]; then
    cp /usr/src/config/ems.yaml /homeassistant
fi

# Copy auth.yaml only if it doesn't already exist
if [ ! -f "/homeassistant/auth.yaml" ]; then
    cp /usr/src/config/auth.yaml /homeassistant
fi


export DENO_INSTALL="/root/.deno"
export PATH="$DENO_INSTALL/bin:$PATH"

deno run -A /usr/src/setupScript.ts /homeassistant/configuration.yaml

bashio::log.info "Done CP JOB"

# Rerun the script every 12 hours to update version info - we don't want to use cron for load reasons
sleep 43200