# ems/rootfs/etc/services.d/add_addon/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Running CP JOB"

mkdir -p /homeassistant/custom_components

cp -R /usr/src/hass-auth-header/custom_components/* /homeassistant/custom_components/

cp /usr/src/config/ems.yaml /homeassistant
cp /usr/src/config/auth.yaml /homeassistant

bashio::log.info "Done CP JOB"

sleep 60000000