# hapro/rootfs/etc/services.d/caddy/run
#!/command/with-contenv bashio
bashio::log.info "Starting Caddy"

source /usr/src/request.sh
doSupervisorRequest "core/info"
export HA_PORT=$(echo "$response" | jq -r .data.port) >> ./env.txt2
export HA_SUBDOMAIN=appel >> ./env.txt2

# doHaProApiRequest "api/Device"
# echo $response

caddy run --config ./Caddyfile --envfile ./env.txt2